# openvas/files/openvas-scanner.service
[Service]
ExecStart=
ExecStart=/usr/sbin/openvassd --listen={{ scanner_default.scanner_address }} --port={{ scanner_default.scanner_port }}
